# distribution/lib/systemd/system/cloud-netconfigd.service
# SPDX-License-Identifier: LGPL-3.0-or-later

[Unit]
Description=Cloud Network Configuration Daemon
Documentation=https://github.com/your-org/cloud-netconfig
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/bin/cloud-netconfigd
Restart=on-failure
RestartSec=5s

# Watchdog
WatchdogSec=60s
NotifyAccess=main

# Security hardening
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectProc=invisible
RestrictAddressFamilies=AF_NETLINK AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallErrorNumber=EPERM

# Runtime directory
RuntimeDirectory=cloud-network
RuntimeDirectoryMode=0755
ConfigurationDirectory=cloud-network
ConfigurationDirectoryMode=0755

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
